# build/rpm/hello.service
[Service]
ExecStart=/usr/bin/hello
Restart=always

[Install]
WantedBy=multi-user.target
